# naa6-audio-bridge.service
[Unit]
Description=NAA6 Audio Bridge — forwards ALSA loopback audio to HQPlayer via NAA 6
Documentation=https://github.com/roon-naa6-bridge
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
ExecStart=/usr/local/bin/naa6-audio-bridge
Restart=on-failure
RestartSec=5s
KillMode=mixed
TimeoutStopSec=10

# Run as a dedicated user for security
User=roon-bridge
Group=audio

# Config and runtime directories
Environment=RUST_LOG=info
RuntimeDirectory=roon-naa6-bridge
RuntimeDirectoryMode=0755

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=naa6-audio-bridge

[Install]
WantedBy=multi-user.target
