# roon-metadata-extension.service
[Unit]
Description=Roon Metadata Extension — forwards track metadata to HQPlayer via NAA 6
Documentation=https://github.com/roon-naa6-bridge
After=network.target naa6-audio-bridge.service
Wants=naa6-audio-bridge.service

[Service]
Type=simple
ExecStart=/usr/bin/node /usr/local/lib/roon-metadata-extension/dist/index.js
Restart=on-failure
RestartSec=5s
KillMode=mixed
TimeoutStopSec=10

# Run as a dedicated user for security
User=roon-bridge
Group=roon-bridge

# Config and runtime directories
Environment=NODE_ENV=production
RuntimeDirectory=roon-naa6-bridge
RuntimeDirectoryMode=0755

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=roon-metadata-extension

[Install]
WantedBy=multi-user.target
